# .config/supervisor/supervisord.conf
[unix_http_server]
file=/tmp/bryan-supervisor.sock   ; the path to the socket file

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:8002        ; ip_address:port specifier, *:port for all iface

[supervisord]
logfile=/dev/null
logfile_backups=0                      ; # of main logfile backups; 0 means none
logfile_maxbytes=0                     ; max main logfile bytes b4 rotation; default 50MB
loglevel=info                          ; log level; default info; others: debug,warn,trace
minfds=1024                            ; min. avail startup file descriptors; default 1024
minprocs=200                           ; min. avail process descriptors;default 200
nodaemon=false                         ; start in foreground if true; default false
pidfile=/tmp/bryan-supervisord.pid     ; supervisord pidfile; default supervisord.pid
silent=false                           ; no logs to stdout if true; default false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/bryan-supervisor.sock ; use a unix:// URL  for a unix socket

[eventlistener:multivisor-rpc]
command=multivisor-rpc --bind *:8052
events=PROCESS_STATE,SUPERVISOR_STATE_CHANGE

[program:clipster]
command=bash -c "wait_for_xserver && clipster -d --config %(ENV_HOME)s/.config/clipster"
redirect_stderr=true
stdout_logfile=%(ENV_LOGDIR)s/%(program_name)s.log
autorestart=true

[group:tools]
programs=clipster

[include]
files=%(here)s/%(host_node_name)s.conf
