# tmp/etc/profile.d/private.sh
export ATHENAS_DDNS_HOSTNAME=bbugyi.ddns.net
export ATHENAS_SSH_PORT=34857

export RFSERVER_HOSTNAME="${ATHENAS_DDNS_HOSTNAME}"
export RFSERVER_TOKEN="P+UDE9+ZCD]vrK8}pdP="
export RFSERVER_PORT=23401
